# systemd/system/package-updates.service
[Unit]
Description=Pacman package updates service
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/package-updates.sh -s
