# systemd/system/package-updates.timer
[Unit]
Description=Pacman package updates timer

[Timer]
OnUnitActiveSec=1h
OnBootSec=10s

[Install]
WantedBy=timers.target
